# vim-manager/vim/templates/skeleton.timer
[Unit]
Description=

[Timer]
OnCalendar=

[Install]
WantedBy=timers.target
